# data/systemd/system/yggdrasil.service.in
[Unit]
Description=yggdrasil system service
Documentation=https://github.com/RedHatInsights/yggdrasil
After=network-online.target
Requires=network-online.target

[Service]
Type=dbus
BusName=com.redhat.Yggdrasil1
NotifyAccess=main
WatchdogSec=300
ExecStart=@bindir@/yggd
PrivateTmp=true

[Install]
WantedBy=multi-user.target
