# data/systemd/system/yggdrasil-bus@.socket
[Unit]
Description=yggdrasil Dispatcher Bus Socket

[Socket]
ListenStream=@yggd_%i

[Install]
WantedBy=sockets.target
