# data/systemd/yggd.service.in
[Unit]
Description=@BRANDNAME@ daemon
Documentation=https://github.com/redhatinsights/yggdrasil
After=network-online.target
Requires=network-online.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=@SBINDIR@/@SHORTNAME@d

[Install]
WantedBy=multi-user.target
