# data/systemd/yggd.service.in
[Unit]
Description=@BRANDNAME@ daemon
Documentation=https://github.com/redhatinsights/yggdrasil

[Service]
Type=simple
ExecStart=@SBINDIR@/@SHORTNAME@d

[Install]
WantedBy=multi-user.target
