# data/systemd/yggd@.service.in
[Unit]
Description=yggdrasil socket service (%i)
Documentation=https://github.com/redhatinsights/yggdrasil
After=network-online.target
Requires=network-online.target yggdrasil-bus@%i.service

[Service]
Type=dbus
BusName=com.redhat.yggdrasil.Dispatcher1
NotifyAccess=main
Environment=DBUS_STARTER_BUS_TYPE=session
Environment=DBUS_SESSION_BUS_ADDRESS=unix:abstract=yggd_%i
ExecStart=@sbindir@/yggd --config @configdir@/yggdrasil-%i.toml

[Install]
WantedBy=multi-user.target
