# data/systemd/yggd.service.in
[Unit]
Description=@SUMMARY@
Documentation=https://github.com/redhatinsights/yggdrasil

[Service]
Type=dbus
BusName=com.redhat.@LONGNAME@1
ExecStart=@SBINDIR@/@SHORTNAME@d

[Install]
WantedBy=multi-user.target
