# data/systemd/yggdrasil@.service.in
[Unit]
Description=yggdrasil socket service (%i)
Documentation=https://github.com/RedHatInsights/yggdrasil
After=network-online.target
Requires=network-online.target yggdrasil-bus@%i.service

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=300
Environment=DBUS_SESSION_BUS_ADDRESS=unix:abstract=yggd_%i
ExecStart=@bindir@/yggd --config @configdir@/yggdrasil-%i.toml
PrivateTmp=true

[Install]
WantedBy=multi-user.target
